# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:r2t-storm]
command=/bin/bash -c "python server.py"
directory=/opt/real-2-time-storm/bin/
autorestart=true
